# ubifs_26/usr/lib/opkg/info/avahi-nodbus-daemon.list
/etc/avahi/http.service
/usr/sbin/avahi-daemon
/etc/config/bonjour
/etc/init.d/avahi-daemon
/etc/avahi/avahi-daemon.conf
/etc/init.d/bonjour
/etc/avahi/https.service
/etc/avahi/avahi-daemon.default_conf
