# ubifs_26/usr/lib/opkg/info/cmm.list
/usr/lib/libcmm.so.0.0.0
/etc/init.d/cmm
/usr/bin/libcmm_sample
/bin/add_vlans
/usr/lib/libcmm.la
/usr/bin/cmm
/usr/lib/libcmm.so.0
/usr/lib/libcmm.so
/usr/lib/libcmm.a
/etc/config/fastforward
