# ubifs_26/usr/lib/opkg/info/cportal.list
/usr/sbin/cportald
/etc/init.d/captive_portal
/usr/sbin/cportalctl
