# ubifs_26/usr/lib/opkg/info/dnsmasq-dhcpv6.list
/etc/dnsmasq.conf
/etc/config/dhcp
/etc/init.d/dnsmasq
/etc/hotplug.d/iface/80-dnsmasq
/usr/sbin/dnsmasq
