# ubifs_26/usr/lib/opkg/info/fci.list
/etc/init.d/fci
/usr/lib/libfci.so
/usr/lib/libfci.so.0
/usr/lib/libfci.so.0.0.0
