# ubifs_26/usr/lib/opkg/info/ipsec-tools.list
/etc/init.d/ipsec
/etc/config/ipsec_new
/usr/sbin/setkey
/etc/racoon/psk.txt
/usr/lib/libracoon.so.0
/etc/racoon.conf
/usr/lib/libracoon.so.0.0.0
/usr/sbin/racoon
/etc/init.d/racoon
/usr/lib/libipsec.so.0
/etc/config/ipsec
/usr/lib/libipsec.so.0.0.1
/usr/sbin/plainrsa-gen
/usr/sbin/racoonctl
