# ubifs_26/usr/lib/opkg/info/kmod-c2k_gpio_reset.list
/etc/init.d/gpio_reset
/lib/modules/4.1.8/c2krv340_gpio_reset.ko
