# ubifs_26/usr/lib/opkg/info/libldap.list
/lib/security/exports.solaris
/etc/ldap.conf
/lib/security/pam_ldap.so
/etc/init.d/ldapclient
/etc/config/ldap
