# ubifs_26/usr/lib/opkg/info/libradius.list
/etc/init.d/radius
/etc/config/radius
/etc/raddb/server
/lib/security/pam_radius_auth.so
