# ubifs_26/usr/lib/opkg/info/netifd.list
/lib/netifd/proto/dhcp.sh
/sbin/pppoe-v6-update
/etc/init.d/network
/lib/netifd/netifd-proto.sh
/sbin/link-hotplug
/sbin/check-pppoe-status
/etc/init.d/netifd
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/lib/network/config.sh
/etc/config/wanport
/etc/hotplug.d/iface/10-sysctl
/sbin/netifd
/sbin/devstatus
/sbin/wan-port-workaround
/sbin/mtu-mac-enable
/lib/netifd/utils.sh
/etc/init.d/wanport
/lib/netifd/netifd-wireless.sh
/sbin/netifd-sync
/sbin/ifstatus
/sbin/ifup
/sbin/route-del
