# ubifs_26/usr/lib/opkg/info/olh.list
/usr/share/olh/RV345-English.tar.gz
/usr/share/olh/RV340W-English.tar.gz
/usr/share/olh/RV345P-English.tar.gz
/usr/share/olh/RV340-English.tar.gz
/etc/init.d/olh
