# ubifs_26/usr/lib/opkg/info/pam-krb5.list
/etc/config/ad
/lib/security/pam_krb5.so
/etc/krb5keytab
/etc/init.d/ad
