# ubifs_26/usr/lib/opkg/info/strongswan.list
/etc/strongswan.d/starter.conf
/usr/sbin/confStrongswanEapRadius
/usr/sbin/timerHandler
/etc/strongswan.d/swanctl.conf
/etc/strongswan.d/pki.conf
/etc/hotplug.d/iface/50-strongswan
/etc/init.d/strongswan
/etc/hotplug.d/iface/13-greUpdate
/usr/lib/ipsec/libstrongswan.so.0
/usr/lib/ipsec/libradius.so
/usr/lib/ipsec/libstrongswan.so.0.0.0
/usr/sbin/tvpnc_timer
/etc/strongswan.d/charon/ezvpn.conf
/etc/ipsec.conf
/etc/config/strongswan
/usr/sbin/check_peer_alive
/usr/sbin/ipxfrm
/usr/lib/ipsec/libvici.so.0.0.0
/usr/sbin/vpnLed
/etc/ipsec.secrets
/usr/sbin/isakmp.bin
/etc/strongswan.conf
/usr/lib/ipsec/libstrongswan.so
/usr/lib/ipsec/plugins/libstrongswan-ezvpn.la
/usr/lib/ipsec/libcharon.so
/usr/lib/ipsec/libvici.so
/usr/lib/ipsec/libvici.so.0
/lib/upgrade/keep.d/strongswan
/usr/sbin/vpnTimer
/etc/init.d/vpnpassthrough
/usr/lib/ipsec/plugins/libstrongswan-ezvpn.so
/usr/sbin/timerCtrl
/etc/strongswan.d/charon.conf
/etc/hotplug.d/iface/45-vpnrouteUpdate
/etc/config/vpnpassthrough
/etc/ipsec.user
/etc/strongswan.d/scepclient.conf
/usr/sbin/vpnLedCtrl
/etc/strongswan.d/charon-logging.conf
