# ubifs_22/usr/lib/opkg/info/802.1x.list
/etc/rtapd.conf
/usr/bin/rtdot1xd
/usr/bin/get_dot1x_portControl
/etc/init.d/dot1x
