# ubifs_22/usr/lib/opkg/info/base-files.list
/etc/firewall6_rv26x-defaults.txt
/etc/firmware_version_minor
/etc/openwrt_version
/default_uci_config/avc
/etc/banner.failsafe
/etc/resolv.conf
/etc/init.d/done
/etc/config/vap0
/bin/config_generate
/etc/config/vap1
/var
/etc/config/vap2
/etc/config/vap3
/etc/config/vap4
/etc/config/vap5
/etc/config/vap6
/etc/config/vap7
/etc/opkg/keys/de98a2dd1d0f8a07
/lib/preinit/00_preinit.conf
/etc/preinit
/default_uci_config/fstab
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/default_uci_config/fastforward
/default_uci_config/webfilter
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/default_uci_config/swupdate
/lib/preinit/10_indicate_preinit
/default_uci_config/l2tpd
/test_scripts/test_button
/default_uci_config/systeminfo
/lib/functions/uci-defaults-new.sh
/etc/firewall_rv26x-defaults.txt
/default_uci_config/dhcp
/etc/config/system
/default_uci_config/device
/default_uci_config/ipsec
/default_uci_config/certificate
/default_uci_config/ad
/etc/fstab
/etc/init.d/zzinfo_show
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/default_uci_config/wlan_led
/sbin/led.sh
/default_uci_config/bootfail
/default_uci_config/ysmir
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/default_uci_config/bonjour
/etc/shadow
/default_uci_config/ddns
/etc/rc.common
/default_uci_config/switch
/etc/uci-defaults/10_migrate-shadow
/default_uci_config/firewall
/test_scripts/test_led
/lib/preinit/10_indicate_failsafe
/default_uci_config/wificonfig
/init
/default_uci_config/dm
/test_scripts/test_usb
/etc/wifi/opt_rv160w.sh
/test_scripts/bridge_all
/default_uci_config/systemstats
/etc/init.d/ucifile_check
/default_uci_config/syslog
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/default_uci_config/crontabs
/sbin/nvram.txt.bak
/default_uci_config/pppoerelay
/default_uci_config/email
/etc/init.d/led
/default_uci_config/igmpproxy
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/init.d/usb
/etc/wifi/pp_rv160w.sh
/etc/firmware_version
/default_uci_config/mwan3
/etc/firewall6_rv16x-defaults.txt
/etc/rc.button/failsafe
/default_uci_config/wireless_old
/default_uci_config/wireless
/default_uci_config/firmware
/bin/check_ubifs_error
/etc/services
/default_uci_config/pptpd
/default_uci_config/wanport
/default_uci_config/vpnpassthrough
/default_uci_config/bwmgmt
/default_uci_config/radvd
/lib/functions/system.sh
/bin/board_detect
/etc/firewall6-defaults.txt
/lib/upgrade/common.sh
/etc/wifi/opt_rv260w.sh
/bin/uci_cache
/etc/sysupgrade.conf
/etc/rc.button/power
/default_uci_config/system
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/init.d/config_update.sh.rv340
/lib/preinit/40_run_failsafe_hook
/etc/group
/default_uci_config/lldpd
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/default_uci_config/ldap
/etc/functions.sh
/default_uci_config/sslvpn
/default_uci_config/ipgroup
/default_uci_config/rip
/default_uci_config/coovachilli
/sbin/firstboot
/etc/swupdateinfo.rv340mr1
/etc/wifi/pp_rv260w.sh
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/etc/hotplug.d/usb/10-usb
/default_uci_config/ipsec_new
/lib/functions.sh
/etc/udhcpc.user
/sbin/power_up.sh
/etc/uci-defaults/11_migrate-sysctl
/etc/config/network
/etc/firewall_rv16x-defaults.txt
/lib/functions/leds.sh
/etc/init.d/zburn_in_test
/rom/note
/etc/inittab
/etc/openwrt_release
/lib/preinit/05_mount_skip
/default_uci_config/strongswan
/sbin/power_down.sh
/sbin/wifi
/default_uci_config/snmpconf
/default_uci_config/upnpd
/etc/profile
/default_uci_config/qos
/default_uci_config/radius
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/default_uci_config/vap0
/default_uci_config/vap1
/default_uci_config/vap2
/default_uci_config/vap3
/default_uci_config/vap4
/default_uci_config/vap5
/default_uci_config/vap6
/default_uci_config/vap7
/etc/init.d/umount
/lib/functions/service.sh
/default_uci_config/ciscosb-yang
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/lib/preinit/02_default_set_state
/default_uci_config/qos_switch
/etc/mtab
/default_uci_config/network
/etc/init.d/boot
/default_uci_config/systemconfig
/lib/functions/preinit.sh
/default_uci_config/poeconf
/etc/firewall-defaults.txt
/default_uci_config/dmz
/bin/ipcalc.sh
/etc/init.d/devices
/etc/diag.sh
/default_uci_config/schedule
/default_uci_config/application
/etc/init.d/config_update.sh
