# ubifs_22/usr/lib/opkg/info/broadcom_esdk_wlan.list
/bin/wps_cli
/usr/sbin/wl
/usr/sbin/acs_cli
/etc/config/cpprofile
/bin/wps_ap_pin
/usr/bin/wps_get_profile.sh
/usr/sbin/nas
/usr/bin/wifi_schedule_check
/usr/bin/wifi_schedule
/usr/bin/wps_nvram.sh
/lib/modules/4.1.8/extra/wl/wl.ko
/usr/bin/wifi_update_affinity
/etc/init.d/wlan_led
/usr/lib/libbcm.so
/bin/wps_daemon
/sbin/wifi_update
/usr/bin/wifi_bttn_check.sh
/bin/wps_pingen
/usr/bin/wifi_set
/etc/config/wlan_led
/usr/bin/wifi_bttn_init.sh
/bin/eapd
/etc/init.d/wifi_switch
/usr/bin/wps_reload.sh
/bin/wps_pbc
/usr/lib/libnvram.so
/sbin/nvram.txt
/usr/lib/libshared.so
/usr/bin/wifi_daemon
/etc/config/wireless
/bin/wps_monitor
/usr/bin/wps_profile.sh
/bin/wps_pin
/usr/bin/initial_ssid_flush_rules.sh
/etc/config/wireless_default
/etc/config/lobby
/usr/lib/libbcmupnp.so
/usr/bin/board_check.sh
/usr/sbin/nvram
/etc/config/wificonfig
/usr/sbin/acsd
/usr/bin/initial_ssid_add_rules.sh
/usr/lib/libbcmcrypto.so
/bin/wps_pin_check
/usr/bin/update_ssid_profile.sh
/etc/config/wlan_led_rv160w
/lib/wifi/broadcom_esdk_wlan.sh
