# ubifs_22/usr/lib/opkg/info/confd.list
/usr/lib/confd/erts/bin/heart
/usr/lib/confd/lib/webui/public_key/include/OTP-PUB-KEY.hrl
/usr/lib/confd/lib/confdc/yanger/priv/yanger_sn.beam
/usr/lib/confd/lib/snmp/snmp/include/SNMP-USM-AES-MIB.hrl
/usr/lib/confd/lib/core/util/priv/pregex_nif.so
/usr/lib/confd/lib/confdc/yanger/priv/yanger_restconf.beam
/usr/lib/confd/lib/core/xds.ez
/usr/lib/confd/lib/confdc/confdc.ez
/usr/lib/confd/lib/webui/yaws/priv/lib/yaws_sendfile_drv.so
/usr/bin/cron.d/5minutes/10-ConfdMonitor.sh
/usr/lib/confd/lib/webui/public_key/include/public_key.hrl
/usr/lib/confd/lib/core/kernel/include/file.hrl
/usr/lib/confd/lib/snmp/snmp/include/STANDARD-MIB.hrl
/usr/lib/confd/lib/snmp/snmpm.ez
/usr/lib/confd/lib/snmp/snmp/include/SNMPv2-MIB.hrl
/usr/lib/confd/bin/confd.boot
/usr/lib/confd/bin/confdc.boot
/usr/lib/confd/lib/core/ssh.ez
/usr/lib/confd/lib/netconf/netconf.ez
/usr/lib/confd/lib/netconf/netconf/priv/yang.fxs
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMPv2-MIB.bin
/usr/lib/confd/lib/core/tts/priv/tts_nif.so
/usr/lib/confd/lib/netconf/netconf/priv/netconf_netmod.fxs
/usr/bin/confd_cli
/usr/lib/confd/lib/core/capi/priv/confd_aaa_bridge
/usr/lib/confd/lib/core/cdb.ez
/usr/lib/confd/lib/core/crypto/priv/lib/crypto.so
/usr/lib/confd/lib/webui/yaws/include/soap.hrl
/usr/bin/confd_cmd
/usr/lib/confd/lib/core/util/priv/syslog_nif.so
/usr/lib/confd/erts/bin/child_setup
/usr/lib/confd/lib/webui/ssl.ez
/usr/lib/confd/lib/snmp/snmp/include/snmp_types.hrl
/usr/lib/confd/lib/core/inets/priv/bin/runcgi.sh
/usr/lib/confd/lib/confdc/yanger/priv/yanger_yang.beam
/usr/lib/confd/lib/core/stdlib/include/zip.hrl
/usr/lib/confd/lib/core/stdlib/include/erl_bits.hrl
/usr/lib/confd/lib/webui/rest/priv/ietf-yang-patch.fxs
/usr/lib/confd/lib/netconf/netconf/priv/netconf_forward.fxs
/usr/bin/confd_load
/usr/lib/confd/lib/confdc/yanger.ez
/usr/lib/confd/lib/webui/yaws/priv/Makefile.in
/etc/confd/stop.d/50_opsdb
/usr/lib/confd/lib/webui/rest/priv/.gitkeep
/usr/lib/confd/lib/confdc/yanger/priv/yanger_ncs.beam
/etc/confd/phase2.d/50_opsdb
/usr/lib/confd/lib/confdc/yanger/bin/yanger
/usr/lib/confd/lib/core/util.ez
/usr/lib/confd/lib/confdc/yanger/priv/yang_parser_nif.so
/usr/lib/confd/erts/bin/escript
/usr/lib/confd/lib/webui/yaws/include/Makefile.am
/usr/lib/confd/lib/core/confd/priv/netconf.fxs
/usr/lib/confd/lib/core/inets/include/httpd.hrl
/usr/lib/confd/lib/confdc/confdc/priv/miba.fxs
/usr/lib/confd/lib/core/util/priv/syst_nif.so
/usr/lib/confd/lib/webui/rest/priv/tailf-rest-error.fxs
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-NOTIFICATION-MIB.bin
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-FRAMEWORK-MIB.bin
/usr/lib/confd/lib/core/util/priv/.gitkeep
/etc/confd/boot.d/00_default
/usr/lib/confd/lib/core/stdlib/include/qlc.hrl
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-TARGET-MIB.bin
/usr/lib/confd/lib/snmp/snmp/include/OTP-SNMPEA-MIB.hrl
/usr/lib/confd/lib/snmp/snmp/include/SNMP-MPD-MIB.hrl
/usr/lib/confd/lib/core/asn1/priv/lib/asn1rt_nif.so
/usr/lib/confd/lib/core/confd/priv/cmdwrapper
/usr/lib/confd/lib/snmp/snmp.ez
/usr/lib/confd/lib/core/asn1.ez
/usr/lib/confd/lib/confdc/yanger/priv/yanger_tailf.beam
/usr/lib/confd/lib/core/compiler.ez
/usr/lib/confd/lib/webui/yaws/include/yaws_dav.hrl
/usr/lib/confd/lib/snmp/snmp/include/SNMPv2-TC.hrl
/usr/bin/confd
/usr/lib/confd/lib/confdc/yanger/priv/.gitkeep
/usr/lib/confd/lib/snmp/snmp/include/SNMP-NOTIFICATION-MIB.hrl
/usr/lib/confd/lib/core/pam/priv/epam
/etc/confd/phase2.d/00_default
/usr/lib/confd/lib/core/confd/priv/crypt-hash
/usr/lib/confd/lib/webui/public_key/include/PKCS-FRAME.hrl
/usr/lib/confd/lib/snmp/snmpgw.ez
/usr/lib/confd/erts/bin/confd.smp
/usr/lib/confd/lib/confdc/yanger/priv/yanger_fxs.beam
/usr/lib/confd/lib/core/tts.ez
/usr/lib/confd/lib/confdc/confdc/priv/clispec.fxs
/usr/lib/confd/lib/netconf/netconf/priv/tailf-netconf-extensions.fxs
/usr/lib/confd/lib/webui/webui.ez
/usr/bin/confd-rpc
/usr/lib/confd/lib/core/kernel/include/inet.hrl
/usr/lib/confd/lib/core/capi.ez
/usr/lib/confd/lib/core/stdlib/include/assert.hrl
/usr/lib/confd/lib/netconf/netconf/priv/netconf_inactive.fxs
/usr/lib/confd/lib/confdc/yanger/priv/yanger_expand.beam
/usr/lib/confd/lib/core/tts/priv/.gitkeep
/usr/lib/confd/lib/netconf/netconf/priv/tailf_netconf.fxs
/usr/lib/confd/lib/core/inets.ez
/usr/lib/confd/lib/core/xpath.ez
/usr/lib/confd/lib/netconf/netconf/priv/netconf_tailf_commit.fxs
/usr/lib/confd/lib/snmp/snmp/priv/mibs/TRANSPORT-ADDRESS-MIB.bin
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-USER-BASED-SM-MIB.bin
/usr/lib/confd/lib/webui/yaws/include/Makefile.in
/usr/lib/confd/lib/snmp/snmp/include/SNMP-VIEW-BASED-ACM-MIB.hrl
/usr/lib/confd/lib/confdc/confdc/priv/index.html.src
/usr/lib/confd/lib/webui/yaws/priv/lib/setuid_drv.so
/usr/lib/confd/lib/netconf/netconf/priv/ietf_with_defaults.fxs
/usr/lib/confd/lib/webui/yaws/include/yaws_api.hrl
/usr/lib/confd/lib/snmp/snmp/include/INET-ADDRESS-MIB.hrl
/usr/lib/confd/lib/core/util/priv/w3cregex_nif.so
/usr/lib/confd/lib/core/stdlib/include/ms_transform.hrl
/etc/confd/phase2.d/99_pnp
/usr/lib/confd/lib/cli/cli/priv/cli-builtin.fxs
/usr/lib/confd/lib/snmp/snmp/include/SNMP-FRAMEWORK-MIB.hrl
/etc/confd/stop.d/00_default
/usr/lib/confd/lib/netconf/netconf/priv/netconf_actions.fxs
/usr/lib/confd/lib/core/crypto.ez
/usr/lib/confd/lib/netconf/netconf/priv/netconf_partial_lock.fxs
/etc/confd/phase0.d/00_default
/etc/confd/phase2.d/60_notification
/usr/lib/confd/lib/snmp/snmp/include/TRANSPORT-ADDRESS-MIB.hrl
/usr/lib/confd/lib/core/confd/priv/ietf-yang-types.fxs
/usr/lib/confd/lib/core/confd.ez
/usr/lib/confd/lib/webui/rest.ez
/usr/lib/confd/lib/core/kernel/include/inet_sctp.hrl
/usr/lib/confd/lib/core/kernel/include/net_address.hrl
/usr/lib/confd/lib/core/kernel/include/dist.hrl
/usr/sbin/config_mgmt.sh
/usr/lib/confd/lib/snmp/snmp/include/SNMPv2-TM.hrl
/usr/lib/confd/lib/confdc/yanger/priv/yanger_tree.beam
/usr/lib/confd/lib/webui/yaws/include/yaws_soap.hrl
/usr/lib/confd/lib/core/capi/priv/utype.so
/usr/lib/confd/lib/core/inets/include/mod_auth.hrl
/usr/lib/confd/lib/webui/yaws/include/yaws.hrl
/usr/lib/confd/erts/bin/confdexec
/usr/lib/confd/lib/webui/yaws/include/soap-envelope.hrl
/usr/lib/confd/bin/erl
/usr/lib/confd/lib/core/msc.ez
/usr/lib/confd/lib/cli/cli.ez
/usr/lib/confd/lib/webui/yaws/priv/mime.types
/usr/lib/confd/lib/snmp/snmp/include/snmp_tables.hrl
/usr/lib/confd/lib/webui/yaws.ez
/usr/lib/confd/lib/netconf/netconf/priv/.gitkeep
/etc/confd/init.d/00_default
/usr/bin/cron.d/hourly/1-config-mirror-update
/usr/lib/confd/lib/snmp/snmp/include/SNMP-TARGET-MIB.hrl
/usr/lib/confd/lib/core/crypto/priv/lib/crypto_callback.so
/usr/lib/confd/lib/core/confd/priv/mmap_schema
/usr/lib/confd/lib/core/confd/priv/confd.fxs
/usr/lib/confd/bin/start_clean.boot
/usr/lib/confd/lib/core/econfd.ez
/usr/lib/confd/erts/bin/confd
/usr/lib/confd/lib/netconf/netconf/priv/netconf_notification.fxs
/usr/lib/confd/bin/escript
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-COMMUNITY-MIB.bin
/usr/lib/confd/lib/webui/yaws/include/yaws_lock.hrl
/usr/lib/confd/lib/snmp/snmp/priv/mibs/IPV6-TC.bin
/usr/lib/confd/lib/confdc/yanger/priv/yanger_smiv2.beam
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-MPD-MIB.bin
/usr/lib/confd/lib/core/kernel/include/dist_util.hrl
/usr/lib/confd/lib/core/stdlib.ez
/usr/lib/confd/lib/snmp/snmp/include/RFC1213-MIB.hrl
/usr/bin/maapi
/usr/lib/confd/lib/core/confd/priv/cmdptywrapper
/usr/lib/confd/lib/core/confd/priv/ipc_drv.so
/usr/lib/confd/lib/core/pam.ez
/usr/lib/confd/lib/cli/cli/priv/.gitkeep
/usr/lib/confd/lib/webui/public_key.ez
/usr/lib/confd/lib/snmp/snmp/priv/mibs/SNMP-VIEW-BASED-ACM-MIB.bin
/usr/lib/confd/lib/webui/rest/priv/tailf-rest-query.fxs
/usr/lib/confd/lib/snmp/snmp/include/SNMP-USER-BASED-SM-MIB.hrl
/etc/init.d/confd
/usr/lib/confd/lib/core/xds/priv/otts_nif.so
/usr/bin/smidump
/usr/lib/confd/lib/core/kernel.ez
/usr/lib/confd/lib/core/confd/priv/config.fxs
/usr/lib/confd/lib/core/confd/priv/confd_cfg.fxs
/usr/lib/confd/erts/bin/inet_gethost
/usr/lib/confd/lib/netconf/ncm.ez
/usr/lib/confd/lib/core/confd/priv/ietf-inet-types.fxs
/usr/lib/confd/lib/webui/yaws/priv/Makefile.am
/usr/bin/confdc
/usr/lib/libconfd.so
/etc/confd/confd.conf
/usr/lib/confd/lib/snmp/snmpa.ez
/usr/lib/confd/lib/core/stdlib/include/erl_compile.hrl
/usr/lib/confd/lib/core/xds/priv/.gitkeep
/usr/lib/confd/lib/confdc/yanger/priv/yanger_lsa.beam
/usr/lib/confd/lib/netconf/netconf/priv/netconf_transactions.fxs
/usr/lib/confd/lib/netconf/netconf/priv/tailf-netconf-query.fxs
/usr/lib/confd/lib/webui/yaws/include/wsdl11soap12.hrl
/etc/confd/phase1.d/00_default
/usr/lib/confd/lib/webui/rest/priv/ietf-restconf.fxs
/usr/lib/confd/bin/erlc
/usr/lib/confd/lib/snmp/snmp/include/SNMP-COMMUNITY-MIB.hrl
