# ubifs_22/usr/lib/opkg/info/epavis-dpi.list
/lib/modules/4.1.8/lionic_dpi.ko
/usr/bin/get-ips-signatures
/root/lionic/sig/application.gz
/usr/bin/lcavc
/etc/config/av
/usr/bin/reset-security-state
/usr/bin/lcstat
/root/lionic/sig/device
/root/lionic/sig/ips.zip.md5
/etc/config/ipsrule
/lib/modules/4.1.8/lionic_pktal.ko
/root/lionic/sig/dm.zip.md5
/usr/bin/get-threats-per-client
/etc/config/dm
/usr/bin/lcsh
/usr/bin/casa
/root/lionic/sig/avc_spec.zip
/root/lionic/sig/ips.zip
/root/lionic/sig/application
/root/lionic/sig/av.zip
/usr/bin/update-clients
/usr/bin/ysmir-sig.sh
/root/lionic/sig/avc.zip.md5
/root/lionic/sig/crf.info
/etc/init.d/avc
/etc/config/ips
/usr/bin/ysmir-start.sh
/usr/bin/decomp_server
/usr/bin/lcips
/usr/bin/uf_agent
/root/lionic/sig/ipsrule
/root/lionic/sig/protocol.tab
/usr/bin/avc-deviceos-metadata.sh
/root/lionic/sig/device.gz
/root/lionic/sig/av.zip.md5
/root/lionic/sig/crf.zip.md5
/root/lionic/sig/mpp.zip
/usr/bin/reset-anti-virus-state
/root/lionic/sig/dm.zip
/root/lionic/sig/ips_spec.zip
/usr/bin/lcav
/root/lionic/sig/crf.siglist.zip
/usr/bin/avc-metadata.sh
/root/lionic/sig/avc.zip
/root/lionic/sig/mpp.zip.md5
/usr/bin/lcjson
/usr/bin/get-attacks-per-client
/usr/bin/lcsig.sh
/lib/modules/4.1.8/lionic_resoft.ko
/etc/config/application
/root/lionic/sig/ipsrule.gz
/root/lionic/sig/bundle.txt
/lib/modules/4.1.8/lionic_fpal.ko
/etc/config/ysmir
/root/lionic/sig/crf.zip
/root/lionic/sig/dm_spec.zip
/usr/bin/ysmir-event.sh
/usr/bin/upnp_cp
/usr/bin/ysmir-config.sh
/usr/bin/get-clients-per-event
/usr/bin/get-top-threats
/lib/modules/4.1.8/lionic_event.ko
/root/lionic/sig/crf.imp2p
/usr/bin/get-top-attacks
/etc/hotplug.d/iface/26-lionic
/root/lionic/sig/dm_oui.zip
/etc/config/device
/usr/bin/ysmir-stop.sh
/root/lionic/sig/dm.tab
/usr/bin/get-top-clients
/usr/bin/reset-ips-state
/etc/config/avc
