# ubifs_22/usr/lib/opkg/info/hardwaredmz.list
/etc/init.d/hardwaredmz
/etc/hotplug.d/iface/94-hardware-dmz
