# ubifs_22/usr/lib/opkg/info/kmod-usb-modem.list
/lib/netifd/proto/directip.sh
/lib/netifd/proto/mobile.sh
/etc/init.d/bwmgmt
/usr/lib/libqmi-glib.so.1.1.0
/etc/gcom/getcnum.gcom
/sbin/usbmode
/etc/usb-mode.json
/etc/gcom/getcarrier.gcom
/etc/modules.d/60-usb-modem
/lib/mobile/data/1199-9051
/lib/mobile/data/1199-9055
/usr/sbin/hilink.sh
/lib/netifd/proto/hilink.sh
/lib/modules/4.1.8/qmi_wwan.ko
/lib/mobile/data/0b05-5602
/etc/gcom/getcardinfo.gcom
/etc/gcom/directip-stop.gcom
/usr/bin/mobile_daemon
/etc/usb-modem.ver
/lib/mobile/mobile-lib.sh
/lib/mobile/data/19d2-1405
/etc/gcom/getcardinfo_ncm.gcom
/lib/mobile/data/12d1-14db
/lib/mobile/data/12d1-14dc
/etc/gcom/3g_signallist
/lib/mobile/data/12d1-1001
/lib/mobile/data/12d1-1003
/etc/gcom/getstrength.gcom
/etc/gcom/getcardinfo_3g.gcom
/usr/bin/cron.d/daily/1-bwmgmt
/lib/modules/4.1.8/usb_wwan.ko
/etc/gcom/getimsi.gcom
/etc/gcom/directip.gcom
/usr/bin/comgt
/sbin/uqmi
/usr/bin/bwmgmt-trigger
/etc/gcom/3g_services
/usr/lib/libqmi-glib.so
/etc/chatscripts/3g.chat
/lib/modules/4.1.8/usbnet.ko
/lib/modules/4.1.8/cdc_ether.ko
/etc/gcom/setmode.gcom
/usr/bin/qmicli
/etc/hotplug.d/tty/30-3g
/usr/sbin/mobile_ifup
/lib/netifd/proto/3g.sh
/lib/netifd/proto/rndis.sh
/lib/modules/4.1.8/cdc-wdm.ko
/etc/chatscripts/evdo.chat
/usr/bin/gcom
/lib/modules/4.1.8/option.ko
/lib/mobile/data/12d1-1506
/etc/hotplug.d/usb/30-mobile
/lib/modules/4.1.8/rndis_host.ko
/lib/modules/4.1.8/huawei_cdc_ncm.ko
/etc/gcom/runcommand.gcom
/lib/netifd/proto/ncm.sh
/etc/gcom/setpin.gcom
/usr/lib/libqmi-glib.so.1
/lib/modules/4.1.8/sierra_net.ko
/usr/bin/qmi-network
/lib/mobile/data/106c-3718
/etc/gcom/ncm.json
/etc/gcom/3g_carrierlist
/usr/bin/bwmgmtd
/etc/gcom/getservicetype.gcom
/usr/bin/reset-mobile-statistics
/lib/mobile/data/1410-9022
/lib/modules/4.1.8/sierra.ko
/lib/mobile/data/1410-9032
/lib/modules/4.1.8/usbserial.ko
/lib/modules/4.1.8/qcserial.ko
/lib/mobile/data/12d1-156c
/lib/mobile/data/0f3d-68aa
/etc/hotplug.d/usb/20-usb_mode
/lib/mobile/data/19d2-0016
/etc/init.d/usbmode
/sbin/usb-modem
/lib/modules/4.1.8/cdc-acm.ko
/lib/mobile/data/19d2-0117
/usr/sbin/mobile_ifdown
/lib/netifd/proto/qmi.sh
