# ubifs_22/usr/lib/opkg/info/krb5-server.list
/usr/sbin/kadmin.local
/usr/bin/sclient
/etc/init.d/krb5kdc
/usr/sbin/krb5kdc
/usr/sbin/kdb5_util
