# ubifs_22/usr/lib/opkg/info/lldpd.list
/etc/config/lldpd
/usr/sbin/lldpd
/etc/init.d/lldpd
/usr/sbin/lldpctl
