# ubifs_22/usr/lib/opkg/info/mcproxy.list
/etc/mcproxy.conf
/etc/config/mcproxy
/etc/hotplug.d/iface/40-mcproxy
/etc/init.d/mcproxy
/usr/sbin/mcproxy
