# ubifs_22/usr/lib/opkg/info/msmtp.list
/etc/msmtprc
/etc/init.d/email
/usr/bin/msmtp
/usr/bin/testserver
/usr/bin/sendmail
/etc/config/email
