# ubifs_22/usr/lib/opkg/info/opsdb.list
/usr/bin/get_wan_port_stats
/usr/bin/update_mobile.sh
/usr/bin/sendOpsdbSignal.sh
/usr/bin/wangetcurrent
/usr/bin/get_upnp_rules
/usr/bin/connect_detect.sh
/usr/bin/l2tpscript
/usr/bin/cps_load_stats
/usr/bin/certscriptCtrl
/usr/bin/get_connections
/usr/bin/wanscript
/usr/bin/aps_load_stats
/usr/bin/neighscriptv4
/usr/bin/neighscriptv6
/lib/liboperdata_api.so
/usr/bin/operdb_stats
/usr/bin/get_lan_port_stats
/usr/bin/get_dhcp_bind_status
/usr/bin/routescript
/usr/bin/lldp_neighscript
/usr/bin/clienttositescript
/etc/services_list
/usr/bin/teleworkerscript
/usr/bin/certscript
/etc/init.d/opsdb
/usr/bin/qosscript
/usr/bin/if_ssid_get_stats
/usr/bin/get_switch_qos_stats
/usr/bin/vap_get_info
/usr/bin/pptpscript
/usr/bin/sitetositescript
/usr/bin/get_wlan_port_stats
/usr/bin/if_radio_get_stats
/usr/bin/ssids_load_stats
/usr/bin/get_port_stats
/usr/bin/executeWanScript
/usr/bin/tcpipservicesscript
