# ubifs_22/usr/lib/opkg/info/pnp_platform.list
/usr/lib/python2.7/site-packages/pnp_platform/discovery/notify.py
/etc/pnp_platform/platform_paths.cfg
/usr/lib/python2.7/site-packages/pnp_platform/infra/platform_constants.py
/usr/lib/python2.7/site-packages/pnp_platform-1.8.0.dev4-py2.7.egg-info/SOURCES.txt
/usr/lib/python2.7/site-packages/pnp_platform/data/bin/create_pnp_config
/usr/lib/python2.7/site-packages/pnp_platform/data/cfg/platform_paths.cfg
/usr/lib/python2.7/site-packages/pnp_platform/services/cli_config.py
/usr/lib/python2.7/site-packages/pnp_platform/services/__init__.py
/etc/init.d/pnpd
/usr/lib/python2.7/site-packages/pnp_platform/infra/__init__.py
/usr/lib/python2.7/site-packages/pnp_platform/data/cfg/platform_logging.cfg
/usr/lib/python2.7/site-packages/pnp_platform/infra/sudi.py
/usr/lib/python2.7/site-packages/pnp_platform/data/init/pnpd
/etc/pnp_platform/policy_discovery.cfg
/usr/lib/python2.7/site-packages/pnp_platform/discovery/lease_parser.py
/usr/lib/python2.7/site-packages/pnp_platform/infra/vrf.py
/usr/lib/python2.7/site-packages/pnp_platform/infra/user_auth.py
/usr/lib/python2.7/site-packages/pnp_platform/infra/device.py
/usr/lib/python2.7/site-packages/pnp_platform/discovery/platform_discovery.py
/usr/lib/python2.7/site-packages/pnp_platform/utils/ntplib.py
/usr/lib/python2.7/site-packages/pnp_platform/discovery/__init__.py
/usr/lib/python2.7/site-packages/pnp_platform-1.8.0.dev4-py2.7.egg-info/top_level.txt
/usr/lib/python2.7/site-packages/pnp_platform/infra/clock_sync.py
/usr/bin/create_pnp_config
/usr/lib/python2.7/site-packages/pnp_platform/utils/__init__.py
/usr/lib/python2.7/site-packages/pnp_platform/infra/event.py
/usr/lib/python2.7/site-packages/pnp_platform/infra/vbond.py
/usr/lib/python2.7/site-packages/pnp_platform/services/image_install.py
/usr/lib/python2.7/site-packages/pnp_platform/services/cli_exec.py
/usr/lib/python2.7/site-packages/pnp_platform/services/config_upgrade.py
/usr/lib/python2.7/site-packages/pnp_platform-1.8.0.dev4-py2.7.egg-info/dependency_links.txt
/usr/lib/python2.7/site-packages/pnp_platform/data/cfg/policy_discovery.cfg
/etc/pnp_platform/platform_logging.cfg
/usr/lib/python2.7/site-packages/pnp_platform/__init__.py
/usr/lib/python2.7/site-packages/pnp_platform-1.8.0.dev4-py2.7.egg-info/PKG-INFO
/usr/lib/python2.7/site-packages/pnp_platform/utils/utilities.py
