# ubifs_22/usr/lib/opkg/info/radvd.list
/usr/sbin/radvd
/etc/hotplug.d/iface/51-radvd
/etc/init.d/radvd
/etc/config/radvd
