# ubifs_22/usr/lib/opkg/info/rdnssd.list
/usr/sbin/rdnssd
/etc/init.d/rdnssd
/etc/merge-hook
