# ubifs_22/usr/lib/opkg/info/rp-pppoe-relay.list
/usr/sbin/pppoe-relay
/etc/init.d/syspppoerelay
/etc/config/pppoerelay
