# ubifs_22/usr/lib/opkg/info/rtk-ms-sdk.list
/usr/sbin/rtkssdk_diag
/usr/sbin/rtk_ms_switch
/etc/init.d/ms-switch
