# ubifs_22/usr/lib/opkg/info/smcroute.list
/usr/bin/mcsender
/usr/sbin/smcroute
/etc/smcroute.conf
/etc/init.d/smcroute
