# ubifs_22/usr/lib/opkg/info/snmpmod.list
/etc/init.d/snmpinit
/usr/bin/getwandns.sh
/usr/bin/getiface.sh
/etc/confd/cdb/snmp_target_init.xml
/usr/bin/snmpglue
/usr/bin/gettrapsconf.sh
/etc/hotplug.d/iface/53-snmp
/usr/bin/snmp_lan_port_stats
/usr/bin/ciscopvtmib_common.sh
/etc/confd/cdb/snmp_init.xml
/etc/confd/cdb/snmp_notify_init.xml
/usr/bin/procdata.sh
/usr/bin/snmpgetcurrent
/usr/bin/snmp_wlan_port_stats
/usr/bin/snmp_lldp_neighinfo.sh
/usr/bin/snmp_wan_port_stats
/etc/confd/cdb/snmp_vacm_init.xml
/usr/bin/snmp_iface_stats.sh
/usr/bin/snmp_port_stats
