# ubifs_22/usr/lib/opkg/info/system.list
/usr/bin/export_diagnostics_PC
/etc/config/schedule
/usr/bin/isSupernet
/usr/bin/delete-upnp-rule
/usr/bin/utility-ping
/etc/init.d/system
/usr/bin/ntpd_update
/usr/bin/certsInfraMigration.sh
/etc/config/systemstats
/usr/bin/cron.d/command/every_week
/usr/bin/utility-dnslookup
/usr/bin/renew-ip
/usr/bin/cron.d/hourly/2-syslog_hourly
/usr/bin/clearconnection_v4.sh
/usr/bin/cron.d/command/every_10minutes
/usr/bin/get_timezoneday
/usr/bin/cron.d/command/every_day
/usr/bin/l2tp_disconnect
/usr/bin/pkt-capture-response
/usr/bin/pptp_disconnect
/usr/bin/reboot_to_image.sh
/usr/bin/pkt-capture-wrapper
/usr/bin/file-delete
/usr/bin/reset-counters
/usr/bin/cron.d/daily/2-syslog_daily
/usr/bin/clock-check-dst
/usr/bin/terminate-cp-user
/usr/bin/reset-avc-state
/usr/bin/reset-qos-statistics
/usr/bin/generate_default_cert
/usr/bin/updateCertUsage
/usr/bin/cap-user-disconnect
/usr/bin/generate-csr
/usr/bin/cron.d/10minutes/2-usage_monitor
/etc/init.d/bootwrapper
/etc/init.d/SwinfoUpdate
/etc/init.d/certificate
/usr/bin/file-upload
/usr/bin/file-copy
/usr/bin/reset-wifi-counters
/usr/bin/updateFlashTime
/usr/bin/clearconnection_v6.sh
/usr/bin/file-download
/usr/bin/connect
/usr/bin/utility-traceroute-request
/usr/bin/system-reboot
/usr/bin/retrieve-wps-pin
/usr/bin/export_diagnostics_USB
/usr/bin/checkDst
/usr/bin/check_usbdevices
/usr/bin/check-latest-version
/usr/bin/getSystemStats
/usr/bin/check_dsttime
/etc/init.d/systemstats
/usr/bin/cron.d/weekly/2-syslog_weekly
/usr/bin/email-server
/etc/mnt_flag
/usr/bin/cron.d/command/every_5minutes
/usr/bin/pkt-capture-request
/usr/bin/c2s-rpc
/etc/config/crontabs
/usr/bin/logread
/etc/config/ipgroup
/usr/bin/s2s-rpc
/usr/bin/bootup_Upgrade
/usr/bin/utility-ping-response
/usr/bin/register-device-pin
/usr/bin/is-port-available
/usr/bin/debug_wireless.sh
/usr/bin/dhcp_release
/usr/bin/dhcp-release-binding
/usr/bin/utility-traceroute
/usr/bin/syserr_codes
/etc/init.d/bootwrapper2
/usr/bin/query-syslog
/usr/bin/generate-cert
/etc/init.d/systemconfig
/etc/config/systeminfo
/usr/bin/cron.d/command/every_month
/mnt/configcert/config/rebootstate
/usr/bin/reset-switch-qos-statistics
/usr/bin/release-ip
/usr/bin/dstHandler
/usr/bin/disconnect
/usr/bin/import-users
/usr/bin/tvpnc-rpc
/usr/bin/cron.d/5minutes/1-healthMonitor.sh
/usr/bin/ifupUsbDevice
/usr/bin/generate-pin
/usr/bin/delete_certificates
/usr/bin/set-current-datetime
/usr/bin/query-syslog-current-page
/usr/bin/install_usb_drivers
/usr/bin/clearconnection.sh
/etc/init.d/lanwancheck
/usr/bin/blink_usb1.sh
/etc/init.d/ipgroup
/etc/config/systemconfig
/etc/init.d/schedule
/usr/bin/utility-ping-request
/usr/bin/cron.d/1minute/1-mwan3OnDemand.sh
/usr/bin/email_diagnostics
/usr/bin/file-list
/usr/bin/cron.d/command/every_hour
/etc/init.d/zfinish
/usr/bin/usb
/etc/boardInfo
/usr/bin/cron.d/command/every_1minute
/usr/bin/press-wps
/usr/bin/blink_usb2.sh
/usr/bin/exportReport.sh
/usr/bin/exec_wanscript.py
