# ubifs_22/usr/lib/opkg/info/webfilter.list
/usr/bin/webfilter-metadata.sh
/usr/bin/url_lookup
/etc/config/webfilter
/usr/bin/wfcli
/usr/bin/wfapp
/etc/init.d/webfilter
/usr/bin/url-lookup
/usr/bin/wf_agent
